# assets/crontab
0 2 * * * cd /opt/yansan-crawler; python -m src.main.crawler.video_info -p