# assets/crontab
0 2 * * * cd /opt/yansan-crawler; python -m crawler.crawler -p