# assets/crontab
0 2 * * * cd /opt/yansan-crawler; python -m src.main.bootstrap.video_info_crawler -p
0 3 * * 1 cd /opt/yansan-crawler; python -m src.main.bootstrap.video_info_cleaner -p
